# codeshrimp/templates/detection/sigma_lotl_attack_detection.yaml
# Purpose: Sigma 检测规则 — Living off the Land (LOTL) 无文件攻击
# Auth: 仅限授权使用 — 仅用于自有/授权网络的 SIEM 部署
# Dependencies: Sigma 2.0+, Winlogbeat/Elastic Agent 日志收集
# Coverage:
#   T1059.001 (PowerShell), T1047 (WMI), T1197 (BITS Jobs),
#   T1218 (Signed Binary Proxy Execution), T1036 (Masquerading),
#   T1055 (Process Injection), T1569 (Service Execution)
# Rule IDs: e5f6a7b8-xxxx-xxx-xxx
# Reference: https://attack.mitre.org/tactics/enterprise/

title: Living off the Land (LOTL) Attack Detection Suite
id: e5f6a7b8-lotl-000-000
status: experimental
author: CodeShrimp
date: 2026/05/11
description: |
  检测利用系统自带工具进行无文件攻击的行为。
  LOTL 攻击不落地恶意文件，使用合法系统工具完成攻击链，
  传统文件扫描无法检测。本规则集覆盖主要 LOTL 技术。

---

# ============================================================
# 1. BITS Job 滥用 — 后台智能传输服务用于下载/持久化
# ============================================================
# T1197 - BITS Jobs
# BITS 可在后台下载文件，绕过网络监控，且支持持久化（重启后续传）

title: Suspicious BITS Job Download from Uncommon URL
id: e5f6a7b8-1197-001-001
status: experimental
description: 检测 BITS 任务从不常见 URL 下载文件，常见于恶意软件下载阶段
author: CodeShrimp
date: 2026/05/11
references:
  - https://attack.mitre.org/techniques/T1197/
tags:
  - attack.persistence
  - attack.defense_evasion
  - attack.t1197
logsource:
  product: windows
  service: bits-client
detection:
  selection:
    EventID: 16384  # BITS job completed
  suspicious_url:
    RemoteName|contains:
      - 'pastebin.com'
      - 'githubusercontent.com'
      - 'mega.nz'
      - 'drive.google.com'
      - 'dropbox.com'
      - '.onion'
      - 'transfer.sh'
      - 'gofile.io'
  suspicious_extension:
    RemoteName|endswith:
      - '.exe'
      - '.ps1'
      - '.bat'
      - '.vbs'
      - '.dll'
      - '.zip'
  condition: selection and (suspicious_url or suspicious_extension)
falsepositives:
  - 合法软件更新使用 BITS
  - 企业内部分发系统
level: high

---

title: BITS Job Persistence via Custom Notification Command
id: e5f6a7b8-1197-001-002
status: experimental
description: BITS 支持在任务完成后执行命令（SetNotifyCmdLine），攻击者利用此特性实现持久化
author: CodeShrimp
date: 2026/05/11
references:
  - https://attack.mitre.org/techniques/T1197/
tags:
  - attack.persistence
  - attack.t1197
logsource:
  product: windows
  service: bits-client
detection:
  selection:
    EventID: 16384
  notify_cmd:
    NotificationCmdLine|contains:
      - 'powershell'
      - 'cmd.exe'
      - 'wscript'
      - 'cscript'
      - 'mshta'
      - 'rundll32'
  condition: selection and notify_cmd
falsepositives:
  - 系统管理工具使用 BITS 通知功能
level: critical

---

# ============================================================
# 2. WMI 持久化与远程执行 — 无文件攻击核心技术
# ============================================================
# T1047 - Windows Management Instrumentation

title: WMI Remote Command Execution via Invoke-WmiMethod
id: e5f6a7b8-1047-001-001
status: experimental
description: 检测通过 WMI 远程执行命令，常见于横向移动
author: CodeShrimp
date: 2026/05/11
references:
  - https://attack.mitre.org/techniques/T1047/
tags:
  - attack.execution
  - attack.lateral_movement
  - attack.t1047
logsource:
  product: windows
  category: process_creation
detection:
  selection:
    Image|endswith:
      - '\powershell.exe'
      - '\pwsh.exe'
    CommandLine|contains|all:
      - 'Invoke-WmiMethod'
      - 'Create'
  suspicious_payload:
    CommandLine|contains:
      - '-enc'
      - '-EncodedCommand'
      - 'IEX'
      - 'Invoke-Expression'
      - 'DownloadString'
      - 'cmd.exe'
      - 'net user'
      - 'net localgroup'
  condition: selection and suspicious_payload
falsepositives:
  - 企业管理软件使用 WMI 远程管理
  - SCCM/Ansible 自动化脚本
level: high

---

title: WMI Process Call Create Remote Execution
id: e5f6a7b8-1047-001-002
status: experimental
description: 检测 Win32_Process.Create 远程调用，WMI 横向移动经典技术
author: CodeShrimp
date: 2026/05/11
tags:
  - attack.execution
  - attack.lateral_movement
  - attack.t1047
logsource:
  product: windows
  category: process_creation
detection:
  selection_wmiprvse:
    ParentImage|endswith: '\wmiprvse.exe'
  suspicious_child:
    Image|endswith:
      - '\cmd.exe'
      - '\powershell.exe'
      - '\pwsh.exe'
  suspicious_cmd:
    CommandLine|contains:
      - '/c'
      - '-enc'
      - 'IEX'
      - 'DownloadString'
      - 'net user'
      - 'reg add'
      - 'schtasks'
  condition: selection_wmiprvse and suspicious_child and suspicious_cmd
falsepositives:
  - 企业管理软件
level: high

---

# ============================================================
# 3. 计划任务滥用 — 持久化与权限提升
# ============================================================
# T1053.005 - Scheduled Task/Job: Scheduled Task

title: Scheduled Task Creation with Suspicious Command
id: e5f6a7b8-1053-005-001
status: experimental
description: 检测创建计划任务时使用可疑命令，如 PowerShell 编码、下载执行等
author: CodeShrimp
date: 2026/05/11
references:
  - https://attack.mitre.org/techniques/T1053/005/
tags:
  - attack.persistence
  - attack.privilege_escalation
  - attack.t1053.005
logsource:
  product: windows
  category: process_creation
detection:
  selection_tool:
    Image|endswith:
      - '\schtasks.exe'
    CommandLine|contains: '/create'
  suspicious_command:
    CommandLine|contains:
      - '-enc'
      - '-EncodedCommand'
      - 'IEX'
      - 'Invoke-Expression'
      - 'DownloadString'
      - 'bypass'
      - '-WindowStyle Hidden'
      - '-NoProfile'
      - 'FromBase64String'
  suspicious_program:
    CommandLine|contains:
      - 'cmd.exe /c'
      - 'wscript'
      - 'cscript'
      - 'mshta'
      - 'rundll32'
      - 'regsvr32'
  condition: selection_tool and (suspicious_command or suspicious_program)
falsepositives:
  - 企业计划任务部署系统
  - 软件安装器创建的更新任务
level: high

---

title: Scheduled Task Created with SYSTEM Privileges by Non-System Process
id: e5f6a7b8-1053-005-002
status: experimental
description: 检测非系统进程以 SYSTEM 权限创建计划任务，可能为权限提升
author: CodeShrimp
date: 2026/05/11
tags:
  - attack.persistence
  - attack.privilege_escalation
  - attack.t1053.005
logsource:
  product: windows
  category: process_creation
detection:
  selection:
    Image|endswith: '\schtasks.exe'
    CommandLine|contains: '/create'
  user_system:
    User: 'NT AUTHORITY\SYSTEM'
  exclude_legitimate:
    ParentImage|endswith:
      - '\services.exe'
      - '\svchost.exe'
      - '\msmpeng.exe'
      - '\tiworker.exe'
      - '\trustedinstaller.exe'
  condition: selection and user_system and not exclude_legitimate
falsepositives:
  - Windows 更新服务
  - 系统维护任务
level: medium

---

# ============================================================
# 4. 服务滥用 — 创建恶意服务/修改现有服务
# ============================================================
# T1569.002 - Service Execution, T1543.003 - Create or Modify System Process: Windows Service

title: Service Created with Suspicious Binary Path
id: e5f6a7b8-1569-002-001
status: experimental
description: 检测创建 Windows 服务时使用可疑二进制路径
author: CodeShrimp
date: 2026/05/11
references:
  - https://attack.mitre.org/techniques/T1569/002/
tags:
  - attack.execution
  - attack.persistence
  - attack.privilege_escalation
  - attack.t1569.002
  - attack.t1543.003
logsource:
  product: windows
  service: system
detection:
  selection:
    EventID: 7045  # Service was installed
  suspicious_path:
    ServiceFileName|contains:
      - 'powershell'
      - 'cmd.exe'
      - 'wscript'
      - 'cscript'
      - 'mshta'
      - 'rundll32'
      - 'regsvr32'
  user_temp:
    ServiceFileName|contains:
      - '\AppData\'
      - '\Temp\'
      - '\Users\Public\'
      - '\ProgramData\'
  suspicious_service_name:
    ServiceName|contains:
      - 'update'
      - 'svchost'
      - 'windows'
      - 'system'
      - 'service'
  condition: selection and (suspicious_path or user_temp or suspicious_service_name)
falsepositives:
  - 合法软件安装
  - 企业服务部署
level: high

---

title: Service Binary Path Modified to Suspicious Value
id: e5f6a7b8-1569-002-002
status: experimental
description: 检测通过 reg/sc 命令修改服务二进制路径为可疑值
author: CodeShrimp
date: 2026/05/11
tags:
  - attack.persistence
  - attack.privilege_escalation
  - attack.t1543.003
logsource:
  product: windows
  category: process_creation
detection:
  selection_reg:
    Image|endswith:
      - '\reg.exe'
    CommandLine|contains:
      - 'add'
      - 'ImagePath'
  selection_sc:
    Image|endswith:
      - '\sc.exe'
    CommandLine|contains:
      - 'config'
      - 'binPath'
      - 'binpath'
  suspicious_value:
    CommandLine|contains:
      - 'powershell'
      - 'cmd.exe'
      - 'IEX'
      - 'DownloadString'
      - 'net user'
      - '\AppData\'
      - '\Temp\'
      - '\Users\Public\'
  condition: (selection_reg or selection_sc) and suspicious_value
falsepositives:
  - 服务配置管理工具
level: high

---

# ============================================================
# 5. Certutil 滥用 — 下载/编码/解码
# ============================================================
# T1105 - Ingress Tool Transfer, T1140 - Deobfuscate/Decode Files

title: Certutil Download or Decode Operation
id: e5f6a7b8-1105-001-001
status: experimental
description: 检测 certutil 用于下载文件（-urlcache）或解码文件（-decode），LOTL 经典技术
author: CodeShrimp
date: 2026/05/11
references:
  - https://attack.mitre.org/techniques/T1105/
  - https://attack.mitre.org/techniques/T1140/
tags:
  - attack.command_and_control
  - attack.defense_evasion
  - attack.t1105
  - attack.t1140
logsource:
  product: windows
  category: process_creation
detection:
  selection:
    Image|endswith: '\certutil.exe'
  download:
    CommandLine|contains:
      - '-urlcache'
      - '-split'
  decode_encode:
    CommandLine|contains:
      - '-decode'
      - '-encode'
  network_target:
    CommandLine|contains:
      - 'http://'
      - 'https://'
      - 'ftp://'
  condition: selection and (download or decode_encode or network_target)
falsepositives:
  - 企业 PKI 管理中的合法 certutil 使用（极少涉及 -urlcache/-decode）
level: high

---

# ============================================================
# 6. MSHTA 滥用 — 执行远程/内联脚本
# ============================================================
# T1218.005 - Signed Binary Proxy Execution: Mshta

title: MSHTA Executing Remote or Inline Script
id: e5f6a7b8-1218-005-001
status: experimental
description: 检测 mshta.exe 执行远程 HTA 或内联 JavaScript/VBScript
author: CodeShrimp
date: 2026/05/11
references:
  - https://attack.mitre.org/techniques/T1218/005/
tags:
  - attack.defense_evasion
  - attack.t1218.005
logsource:
  product: windows
  category: process_creation
detection:
  selection:
    Image|endswith: '\mshta.exe'
  remote_execution:
    CommandLine|contains:
      - 'http://'
      - 'https://'
  inline_script:
    CommandLine|contains:
      - 'javascript:'
      - 'vbscript:'
      - 'eval('
      - 'CreateObject'
      - 'WScript.Shell'
      - 'Shell.Application'
  suspicious_parent:
    ParentImage|endswith:
      - '\cmd.exe'
      - '\powershell.exe'
      - '\wscript.exe'
      - '\cscript.exe'
      - '\winword.exe'
      - '\excel.exe'
  condition: selection and (remote_execution or inline_script or suspicious_parent)
falsepositives:
  - HTA 应用程序正常使用（极少）
  - SCCM 部署脚本
level: high

---

# ============================================================
# 7. 进程注入检测 — 无文件攻击横向移动
# ============================================================
# T1055 - Process Injection

title: Process Injection via Rundll32 with Suspicious DLL
id: e5f6a7b8-1055-001-001
status: experimental
description: 检测 rundll32 加载可疑 DLL，常见于进程注入和代理执行
author: CodeShrimp
date: 2026/05/11
references:
  - https://attack.mitre.org/techniques/T1055/
tags:
  - attack.defense_evasion
  - attack.privilege_escalation
  - attack.t1055
  - attack.t1218.011
logsource:
  product: windows
  category: process_creation
detection:
  selection:
    Image|endswith: '\rundll32.exe'
  suspicious_dll:
    CommandLine|contains:
      - '\AppData\'
      - '\Temp\'
      - '\Users\Public\'
      - '\ProgramData\'
      - 'javascript:'
      - 'mshtml'
  suspicious_entry:
    CommandLine|contains:
      - 'DllRegisterServer'
      - 'DllGetClassObject'
      - 'Start'
      - 'EntryPoint'
  no_parent_explorer:
    ParentImage|endswith:
      - '\explorer.exe'
  condition: selection and (suspicious_dll or suspicious_entry) and not no_parent_explorer
falsepositives:
  - 合法软件使用 rundll32 调用系统 DLL
level: medium

---

title: Suspicious Process Hollowing Indicators
id: e5f6a7b8-1055-001-002
status: experimental
description: 检测进程镂空（Process Hollowing）的命令行特征
author: CodeShrimp
date: 2026/05/11
tags:
  - attack.defense_evasion
  - attack.t1055.012
logsource:
  product: windows
  category: process_creation
detection:
  selection_legit_process:
    Image|endswith:
      - '\svchost.exe'
      - '\csrss.exe'
      - '\winlogon.exe'
      - '\lsass.exe'
      - '\wininit.exe'
      - '\dwm.exe'
      - '\taskhostw.exe'
  suspicious_parent:
    ParentImage|endswith:
      - '\cmd.exe'
      - '\powershell.exe'
      - '\wscript.exe'
      - '\cscript.exe'
      - '\mshta.exe'
  no_commandline:
    CommandLine: ''
  condition: selection_legit_process and suspicious_parent and no_commandline
falsepositives:
  - 正常系统进程管理
level: critical

---

# ============================================================
# 8. 综合场景 — 完整 LOTL 攻击链检测
# ============================================================
# 检测从初始访问到横向移动的完整 LOTL 攻击链

title: LOTL Attack Chain - Multiple LOLBIN Execution in Short Timeframe
id: e5f6a7b8-lotl-chain-001
status: experimental
description: |
  检测短时间内连续使用多个系统工具（LOLBIN），可能为 LOTL 攻击链。
  攻击者通常在短时间内连续使用 certutil(bitsadmin)下载→mshta执行→powershell横向移动。
author: CodeShrimp
date: 2026/05/11
tags:
  - attack.execution
  - attack.defense_evasion
  - attack.lateral_movement
  - attack.t1059.001
  - attack.t1218
  - attack.t1047
logsource:
  product: windows
  category: process_creation
detection:
  selection_lolbin:
    Image|endswith:
      - '\certutil.exe'
      - '\bitsadmin.exe'
      - '\mshta.exe'
      - '\wmic.exe'
      - '\regsvr32.exe'
      - '\rundll32.exe'
      - '\msiexec.exe'
      - '\installutil.exe'
  selection_script:
    Image|endswith:
      - '\powershell.exe'
      - '\pwsh.exe'
      - '\wscript.exe'
      - '\cscript.exe'
  selection_cmd:
    Image|endswith:
      - '\cmd.exe'
  condition: selection_lolbin or selection_script or selection_cmd
  timeframe: 5m
  condition2: count(Image) by SourceImage > 5
falsepositives:
  - 企业软件部署和更新
  - 系统管理自动化脚本
level: medium
